# HOME/.config/autostart/tilda.desktop
[Desktop Entry]
Name=Tilda
Exec=tilda
Icon=tilda
Terminal=false
Type=Application
Categories=GNOME;GTK;Utility;TerminalEmulator;System;
